# files/celery-supervisor.conf
[program:celeryd]
command=/var/www/site/venv/bin/python manage.py celery worker
directory=/var/www/site/
user=nobody
autostart=true
autorestart=true
