# files/site-supervisor.conf
[program:site]
command=/var/www/site/venv/bin/gunicorn myapp.wsgi:application --workers=3
directory=/var/www/site/
user=nobody
autostart=true
autorestart=true
